# backup.service
[Unit]
Description=Minecraft Reset

[Service]
ExecStart=/var/minecraft/backup.sh

[Install]
WantedBy=multi-user.target
