# mcreset.service
[Unit]
Description=Minecraft Reset

[Service]
ExecStart=/var/minecraft/mcreset.sh

[Install]
WantedBy=multi-user.target
