# ramdisk.service
[Unit]
Description=Minecraft Ramdisk

[Service]
ExecStart=/var/minecraft/ramdisk.sh

[Install]
WantedBy=multi-user.target
