# multiserver/mcreset01.service
[Unit]
Description=Minecraft Reset

[Service]
ExecStart=/var/minecraft/mcreset01.sh

[Install]
WantedBy=multi-user.target
