# multiserver/multilaunch.service
[Unit]
Description=Minecraft Reset

[Service]
ExecStart=/var/minecraft/multilaunch.sh

[Install]
WantedBy=multi-user.target
